# install_scripts/rpi4-22.04/rc.local
#!/usr/bin/env bash

# Fan is on pwm5
echo 5 > /sys/class/pwm/pwmchip0/export
PERIOD=20000000
echo $PERIOD > /sys/class/pwm/pwmchip0/pwm5/period
echo $(( ${PERIOD} / 2)) > /sys/class/pwm/pwmchip0/pwm5/duty_cycle